# multi-agent-swarm.service
[Unit]
Description=Multi-Agent Swarm System (CrewAI)
After=network.target

[Service]
Type=simple
User=__YOUR_LINUX_USER__
WorkingDirectory=/opt/multi-agent-swarm
Environment="GCP_PROJECT_ID=__YOUR_GCP_PROJECT_ID__"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/opt/multi-agent-swarm/service-account.json"
ExecStart=/opt/multi-agent-swarm/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/multi-agent-swarm.log
StandardError=append:/var/log/multi-agent-swarm.error.log

[Install]
WantedBy=multi-user.target
